# pingloop.slice
[Unit]
Sescription=Slice for pingloop service

[Slice]
CPUQuota=5%
MemoryLimit=40M
